# mac-sender/com.keybridgebt.sender.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.keybridgebt.sender</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/keybridgebt/mac-sender/.venv/bin/python3</string>
        <string>-m</string>
        <string>keybridgebt_mac</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/keybridgebt/mac-sender</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>UserName</key>
    <string>_keybridgebt</string>

    <key>StandardOutPath</key>
    <string>/var/log/keybridgebt/sender.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/keybridgebt/sender.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
